# target/linux/ath79/base-files/etc/init.d/qca8k_set_port
#!/bin/sh /etc/rc.common

START=15

set_qca8k_port() {
	local port=$1
	local master=$2

	ip link set $port type dsa conduit $master
}

boot() {
	# Restore original implementation where the eth1 (port 6) was used
	# for the lan port and the eth0 (port 0) was used for the wan port
	case $(board_name) in
		tplink,archer-c7-v2)
			set_qca8k_port lan1 eth1
			set_qca8k_port lan2 eth1
			set_qca8k_port lan3 eth1
			set_qca8k_port lan4 eth1
			set_qca8k_port wan eth0
			;;
	esac
}
