# target/linux/ath79/base-files/etc/uci-defaults/01_disable_gro
# SPDX-License-Identifier: GPL-2.0-only

uci set network.device_defaults='device_defaults'
uci set network.device_defaults.gro='false'
